# curator4-771741017.conf
[program:curator4-771741017]
command=python3 /home/curator4/curator4-771741017/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-771741017/test.err.log
stdout_logfile=/home/curator4/curator4-771741017/test.out.log